# hwraid-smartarray/crontab
# Description:	HP SmartArray data gathering and analyze
0 6 * * * /usr/libexec/zabbix-extensions/scripts/hp-raid-data-processor.sh
